# deploy/systemd/convex-local.service
[Unit]
Description=OpenBrain Convex Local Deployment
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/openbrain
ExecStart=/usr/bin/env bash -lc 'cd %h/openbrain && npx convex dev --tail-logs disable'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target

