# deploy/systemd/lmstudio.service
[Unit]
Description=LM Studio headless daemon and API server
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=%h/openbrain
Environment="PATH=%h/.lmstudio/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LMSTUDIO_BIN=%h/.lmstudio/bin/lms"
Environment="LMSTUDIO_MODEL=text-embedding-embeddinggemma-300m-qat"
Environment="LMSTUDIO_RUNTIME_ALIAS=llama.cpp-linux-x86_64-vulkan-avx2"
Environment="LMSTUDIO_GPU_OFFLOAD=max"
Environment="LMSTUDIO_REQUIRE_VULKAN=1"
Environment="LMSTUDIO_BASE_URL=http://127.0.0.1:1234/v1"
Environment="LMSTUDIO_STARTUP_TIMEOUT_SECONDS=45"
Environment="LMSTUDIO_HEALTH_INTERVAL_SECONDS=5"
ExecStartPre=-%h/.lmstudio/bin/lms server stop
ExecStartPre=-%h/.lmstudio/bin/lms daemon down
ExecStart=/usr/bin/env bash -lc 'set -euo pipefail; cd %h/openbrain; exec %h/openbrain/scripts/lmstudio-supervisor.sh'
ExecStop=-%h/.lmstudio/bin/lms server stop
ExecStop=-%h/.lmstudio/bin/lms daemon down
Restart=always
RestartSec=3
TimeoutStartSec=240
TimeoutStopSec=30

[Install]
WantedBy=default.target
