# deploy/systemd/lmstudio.service.tmpl
[Unit]
Description=LM Studio headless daemon and API server
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
Environment="PATH=__USER_HOME__/.lmstudio/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LMSTUDIO_BIN=__LMSTUDIO_BIN__"
Environment="LMSTUDIO_MODEL=__LMSTUDIO_MODEL__"
Environment="LMSTUDIO_RUNTIME_ALIAS=__LMSTUDIO_RUNTIME_ALIAS__"
Environment="LMSTUDIO_GPU_OFFLOAD=__LMSTUDIO_GPU_OFFLOAD__"
Environment="LMSTUDIO_REQUIRE_VULKAN=__LMSTUDIO_REQUIRE_VULKAN__"
Environment="LMSTUDIO_BASE_URL=__LMSTUDIO_BASE_URL__"
Environment="LMSTUDIO_STARTUP_TIMEOUT_SECONDS=__LMSTUDIO_STARTUP_TIMEOUT_SECONDS__"
Environment="LMSTUDIO_HEALTH_INTERVAL_SECONDS=__LMSTUDIO_HEALTH_INTERVAL_SECONDS__"
ExecStartPre=-__LMSTUDIO_BIN__ server stop
ExecStartPre=-__LMSTUDIO_BIN__ daemon down
ExecStart=/usr/bin/env bash -lc 'set -euo pipefail; cd __WORKDIR__; exec __WORKDIR__/scripts/lmstudio-supervisor.sh'
ExecStop=-__LMSTUDIO_BIN__ server stop
ExecStop=-__LMSTUDIO_BIN__ daemon down
Restart=always
RestartSec=3
TimeoutStartSec=240
TimeoutStopSec=30

[Install]
WantedBy=default.target
