# deploy/systemd/lmstudio.service.tmpl
[Unit]
Description=LM Studio headless daemon and API server
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="PATH=__USER_HOME__/.lmstudio/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LMSTUDIO_MODEL=__LMSTUDIO_MODEL__"
Environment="LMSTUDIO_RUNTIME_ALIAS=__LMSTUDIO_RUNTIME_ALIAS__"
Environment="LMSTUDIO_GPU_OFFLOAD=__LMSTUDIO_GPU_OFFLOAD__"
Environment="LMSTUDIO_REQUIRE_VULKAN=__LMSTUDIO_REQUIRE_VULKAN__"
ExecStart=__LMSTUDIO_BIN__ daemon up
ExecStart=__LMSTUDIO_BIN__ runtime select ${LMSTUDIO_RUNTIME_ALIAS}
ExecStart=__LMSTUDIO_BIN__ load ${LMSTUDIO_MODEL} --gpu ${LMSTUDIO_GPU_OFFLOAD} --yes
ExecStart=__LMSTUDIO_BIN__ server start
ExecStartPost=/usr/bin/env bash -lc 'if [ "${LMSTUDIO_REQUIRE_VULKAN}" = "1" ]; then __LMSTUDIO_BIN__ runtime survey | grep -qi "Survey by .*vulkan"; fi'
ExecStop=__LMSTUDIO_BIN__ server stop
ExecStop=__LMSTUDIO_BIN__ daemon down

[Install]
WantedBy=default.target
